# sensor/platforms/bbb/usb_audit/usb-audit@.service
[Service]
Type=oneshot
ExecStart=/opt/usb_audit.sh %i