# sensor/platforms/docker_x86/services/manager/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=$(find /etc/manager -type f -print -mindepth 1 -maxdepth 1 | head -1)
if [ -z "$CONFIG_FILE" ]; then
  echo "Error: No config file found"
  exit 1
else
  echo "Using config file $CONFIG_FILE"
fi

if [ "$DEV_MODE" ]; then
  echo "Starting manager in development mode"
  DEV_MODE="-d"
fi

if [ "$IFACE" ]; then
  IFACE="-i $IFACE"
fi

if [ "$LOG_LVL" ]; then
  LOG_LVL="-l $LOG_LVL"
fi

# Force unbuffered stdout
export PYTHONUNBUFFERED=1

exec manager -p docker $DEV_MODE $IFACE $LOG_LVL $CONFIG_FILE
