# sensor/platforms/docker_x86/services/manager/run
#!/usr/bin/with-contenv sh

CONFIG_FILE=$(find /etc/manager -type f -print -mindepth 1 -maxdepth 1 | head -1)
if [ -z "$CONFIG_FILE" ]; then
  echo "Error: No config file found"
  exit 1
else
  echo "Using config file $CONFIG_FILE"
fi

if [ "$SKIP_INIT" ]; then
  echo "Development mode: Skipping network initialization"
  SKIP_INIT="-s"
fi

IFACE=""
if [ -n "$INTERFACE" ]; then
    IFACE="-i $INTERFACE"
fi

# Force unbuffered stdout
export PYTHONUNBUFFERED=1

exec manager -p docker $SKIP_INIT $IFACE $CONFIG_FILE
