# server/utils/docker/services/update-worker/run
#!/bin/bash

# Force unbuffered stdout
export PYTHONUNBUFFERED=1

exec /sbin/setuser beanstalkd /usr/bin/python /opt/HoneySens/beanstalk/honeysens-update-worker.py /opt/HoneySens/data/config.cfg
