# sensor/platforms/docker_x86/services/cntlm/run
#!/usr/bin/with-contenv sh

exec cntlm -f