# server/utils/docker/services/beanstalkd/run
#!/bin/bash
exec /sbin/setuser beanstalkd /usr/bin/beanstalkd -l 127.0.0.1 -p 11300
